# crontab
0 5  * * *  /home/ubuntu/microsvc-lecrynginx-docker/cron_job.sh
